# server/crontab
# Xarticle Cron Jobs
# Run on Racknerd server

# Article Pipeline - Run every 2 hours
# Fetch tweets, extract articles, generate summaries, save to database
0 */2 * * * cd /root/xarticle-server && npm run pipeline >> logs/pipeline.log 2>&1

# Daily Report - Run every day at 7:00 AM UTC
0 7 * * * cd /root/xarticle-server && npm run daily-report >> logs/daily-report.log 2>&1

# Optional: Run pipeline more frequently during peak hours
# 30 8,10,12,14,16,18,20,22 * * * cd /root/xarticle-server && npm run pipeline >> logs/pipeline.log 2>&1
